# ldap-client/run7.service
[Unit]

[Service]
ExecStart=/root/run7.sh

[Install]
WantedBy=default.target
